# motherbrain/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock 

[supervisord]
logfile=/tmp/supervisord.log 
loglevel=error
logfile_maxbytes=50MB       
logfile_backups=10         
pidfile=/tmp/supervisord.pid
nodaemon=false            
minfds=1024          
minprocs=200      

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:Daphne]
environment=PATH="/opt/python/run/venv/bin"
command=sh /opt/python/current/app/motherbrain/start_daphne.sh --log-file /tmp/start_daphne.log
directory=/opt/python/current/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/daphne.out.log
stderr_logfile=/tmp/daphne.err.log
